# tools/edge_scan_feasibility/service/m20pro-edge-scan-106.service.example
[Unit]
Description=M20Pro 106 production edge scan service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
EnvironmentFile=/etc/m20pro-edge-scan-106.env
ExecStart=/bin/bash -lc 'source /opt/robot/scripts/setup_ros2.sh && exec /usr/local/lib/m20pro/m20pro_edge_scan "${INPUT_TOPIC}" "${OUTPUT_TOPIC}" "${DURATION_S}" "${DOMAIN}" "${USE_SHM}" "${PREFIX}" "${HEIGHT_MIN}" "${HEIGHT_MAX}" "${MAX_PUBLISH_HZ}" "${MAX_POINTS}" "${FRAME_ID}" "${ANGLE_INCREMENT}" "${RANGE_MAX}" "${RANGE_MIN}" "${BIN_HOLD_S}" "${FIELD_PROFILE_NAME}" "${FIELD_PROFILE_HASH}"'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
